# rootfs/etc/crontabs/caddy
0 8 * * * /usr/local/bin/cli cronjob > /dev/null 2>&1
